# root/etc/init.d/hp-lites
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

EXTRA_COMMANDS="write_config clear_log"

PROG=/usr/bin/hp-lites
WORKDIR=/etc/hp-lite
CONFIG_FILE="$WORKDIR/app.yml"
DEFAULT_LOG_DIR=/var/log/hp-lite
LOG_FILE_NAME=hp-lites.log
DEFAULT_LOG_FILE=$DEFAULT_LOG_DIR/$LOG_FILE_NAME
LOG_FILE="$DEFAULT_LOG_FILE"
LOG_DIR="$DEFAULT_LOG_DIR"
CRON_TAG="# hp-lites-log-clean"
DEFAULT_RETENTION=3

shell_quote() {
    printf "'%s'" "$(printf "%s" "$1" | sed "s/'/'\\\\''/g")"
}

uci_get_default() {
    local value
    value="$(uci get "$1" 2>/dev/null)"
    [ -n "$value" ] && printf "%s" "$value" || printf "%s" "$2"
}

yaml_quote() {
    shell_quote "$1"
}

yaml_bool() {
    case "$1" in
        1|true|TRUE|yes|YES|on|ON) printf "true" ;;
        *) printf "false" ;;
    esac
}

normalize_retention() {
    local ret="$1"

    [ -z "$ret" ] && ret=$DEFAULT_RETENTION
    case "$ret" in
        *[!0-9]*)
            ret=$DEFAULT_RETENTION
            ;;
    esac
    [ "$ret" -lt 1 ] && ret=$DEFAULT_RETENTION
    [ "$ret" -gt 3650 ] && ret=$DEFAULT_RETENTION

    printf "%s" "$ret"
}

normalize_log_dir() {
    local dir="$1"

    [ -z "$dir" ] && dir="$DEFAULT_LOG_DIR"
    case "$dir" in
        /*) ;;
        *) dir="$DEFAULT_LOG_DIR" ;;
    esac
    case "$dir" in
        *..*) dir="$DEFAULT_LOG_DIR" ;;
    esac
    while [ "$dir" != "/" ] && [ "${dir%/}" != "$dir" ]; do
        dir="${dir%/}"
    done

    printf "%s" "$dir"
}

normalize_log_file() {
    local file="$1"

    file="$(printf "%s" "$file" | sed 's/^[[:space:]]*//;s/[[:space:]]*$//')"
    [ -z "$file" ] && file="$DEFAULT_LOG_FILE"
    case "$file" in
        /*) ;;
        *) file="$DEFAULT_LOG_FILE" ;;
    esac
    case "$file" in
        *..*|*/) file="$DEFAULT_LOG_FILE" ;;
    esac

    printf "%s" "$file"
}

log_dir_from_file() {
    local file="$1"
    local dir="${file%/*}"

    [ "$dir" = "$file" ] && dir="/"
    [ -z "$dir" ] && dir="/"

    printf "%s" "$dir"
}

set_log_paths() {
    local file
    local dir

    file="$(uci get hp-lites.server.log_file 2>/dev/null)"
    if [ -z "$file" ]; then
        dir="$(normalize_log_dir "$(uci get hp-lites.server.log_dir 2>/dev/null)")"
        file="$dir/$LOG_FILE_NAME"
    fi

    LOG_FILE="$(normalize_log_file "$file")"
    LOG_DIR="$(log_dir_from_file "$LOG_FILE")"
}

schedule_log_cleanup() {
    local ret
    local clean_cmd

    set_log_paths
    ret="$(normalize_retention "$1")"
    clean_cmd="mkdir -p $(shell_quote "$LOG_DIR") && : > $(shell_quote "$LOG_FILE")"

    [ -f /etc/crontabs/root ] || touch /etc/crontabs/root
    sed -i "/$CRON_TAG/d" /etc/crontabs/root
    echo "0 3 */$ret * * $clean_cmd $CRON_TAG" >> /etc/crontabs/root
    /etc/init.d/cron reload >/dev/null 2>&1 || /etc/init.d/cron restart >/dev/null 2>&1
}

write_config() {
    local admin_username
    local admin_password
    local admin_port
    local cmd_port
    local tunnel_ip
    local tunnel_port
    local open_domain
    local acme_email
    local acme_http_port

    mkdir -p "$WORKDIR"

    admin_username="$(uci_get_default hp-lites.server.admin_username zyh)"
    admin_password="$(uci_get_default hp-lites.server.admin_password 123456)"
    admin_port="$(uci_get_default hp-lites.server.admin_port 9090)"
    cmd_port="$(uci_get_default hp-lites.server.cmd_port 16666)"
    tunnel_ip="$(uci_get_default hp-lites.server.tunnel_ip 127.0.0.1)"
    tunnel_port="$(uci_get_default hp-lites.server.tunnel_port 9091)"
    open_domain="$(yaml_bool "$(uci_get_default hp-lites.server.open_domain 0)")"
    acme_email="$(uci_get_default hp-lites.server.acme_email "1540187368@qq.com")"
    acme_http_port="$(uci_get_default hp-lites.server.acme_http_port 5634)"

    cat > "$CONFIG_FILE" <<EOC
admin:
  username: $(yaml_quote "$admin_username")
  password: $(yaml_quote "$admin_password")
  port: $admin_port

cmd:
  port: $cmd_port

tunnel:
  ip: $(yaml_quote "$tunnel_ip")
  port: $tunnel_port
  open-domain: $open_domain

acme:
  email: $(yaml_quote "$acme_email")
  http-port: $(yaml_quote "$acme_http_port")
EOC

    chmod 600 "$CONFIG_FILE"
}

clear_log() {
    set_log_paths
    mkdir -p "$LOG_DIR"
    : > "$LOG_FILE"
    chmod 644 "$LOG_FILE"
}

start_service() {
    set_log_paths
    mkdir -p "$LOG_DIR"
    mkdir -p "$WORKDIR"
    touch "$LOG_FILE"
    chmod 644 "$LOG_FILE"

    echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lites] Starting service..." >> "$LOG_FILE"

    [ -x "$PROG" ] || {
        echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lites] Start failed: $PROG is not installed" >> "$LOG_FILE"
        return 1
    }

    write_config

    local cmd
    cmd="$(shell_quote "$PROG") -conf $(shell_quote "$CONFIG_FILE") 2>&1 | tee -a $(shell_quote "$LOG_FILE")"

    procd_open_instance
    procd_set_param command sh -c "$cmd"
    procd_set_param respawn
    procd_set_param cwd "$WORKDIR"
    procd_set_param file /etc/config/hp-lites
    procd_set_param file "$CONFIG_FILE"
    procd_close_instance

    local RET
    RET=$(uci get hp-lites.server.log_retention_days 2>/dev/null)
    schedule_log_cleanup "$RET"

    echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lites] Service started successfully, process is supervised by procd" >> "$LOG_FILE"
    echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lites] Log cleanup scheduled: clear server log every $(normalize_retention "$RET") days at 03:00" >> "$LOG_FILE"
}

stop_service() {
    set_log_paths
    mkdir -p "$LOG_DIR"
    touch "$LOG_FILE"
    chmod 644 "$LOG_FILE"

    procd_kill hp-lites 2>/dev/null

    local pids
    pids=$(pgrep -f "$PROG")
    if [ -n "$pids" ]; then
        local pid
        for pid in $pids; do
            if kill -9 "$pid" 2>/dev/null; then
                echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lites] Process (PID: $pid) terminated successfully" >> "$LOG_FILE"
            else
                echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lites] Failed to terminate process (PID: $pid)" >> "$LOG_FILE"
            fi
        done

        pkill -9 -f "$PROG" 2>/dev/null
        echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lites] Cleaned up all hp-lites related processes" >> "$LOG_FILE"
    else
        echo "[$(date +"%Y-%m-%d %H:%M:%S")] [hp-lites] No running hp-lites process found" >> "$LOG_FILE"
    fi
}
